# docker/graphite-watcher/supervisord.conf
[supervisord]
nodaemon=true

[program:graphite-watcher]
command=/go/bin/graphite-watcher
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/graphite-watcher.log
