# docker/tsdb-gw/supervisord.conf
[supervisord]
nodaemon=true

[program:tsdb-gw]
command=/go/bin/tsdb-gw
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/tsdb-gw.log
autostart=false
