# docker/eventtank/supervisord.conf
[supervisord]
nodaemon=true

[program:eventtank]
command=/go/bin/eventtank
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/eventtank.log
autostart=false
