# images/apps-server/supervisord.conf
[supervisord]
nodaemon=true

[program:apps-server]
command=/start.sh
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/apps-server.log
startretries=10
