# nsq_events/supervisord.conf
[supervisord]
nodaemon=true

[program:nsq_to_kairos]
command=/go/bin/nsq_events
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/nsq_events.log
