# docker/graphitewatcher/supervisord.conf
[supervisord]
nodaemon=true

[program:graphitewatcher]
command=/go/bin/graphitewatcher
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/graphitewatcher.log
