# kairosdb/supervisord.conf
[supervisord]
nodaemon=true

[program:kairosdb]
command=start_kairosdb.py
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/kairosdb.log
startsecs=10
startretries=5
