# images/metric_tank/supervisord.conf
[supervisord]
nodaemon=true

[program:metric_tank]
command=/start.sh
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/metric_tank.log
