# images/tsdb/supervisord.conf
[supervisord]
nodaemon=true

[program:tsdb]
command=/start.sh
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/tsdb.log
startretries=10
