# images/worldping-api/supervisord.conf
[supervisord]
nodaemon=true

[program:worldping-api]
command=/start.sh
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/worldping-api.log
startretries=10
