# images/task-server/supervisord.conf
[supervisord]
nodaemon=true

[program:task-server]
command=/start.sh
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile=/var/log/raintank/task-server.log
startretries=10
